# web/deploy/researchvault.service
[Unit]
Description=ResearchVault API (rv.insas.jp only)
Documentation=https://github.com/sohiko/ResearchVault
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rv
Group=rv

# API: server.js + api/ + node_modules
WorkingDirectory=/home/rv/www/app
# 静的ファイルは既存デプロイ先（GitHub Actions ビルド → public/）
Environment=DIST_DIR=/home/rv/www/public
Environment=NODE_ENV=production
# dev.insas.jp が 3000 使用のため rv は 3001（insas.jp には触れない）
Environment=PORT=3001

ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=researchvault

[Install]
WantedBy=multi-user.target
